# deploy/lagersystem.service
[Unit]
Description=Lagersystem Next.js app
After=network.target

[Service]
Type=simple
User=mimer
Group=mimer
WorkingDirectory=/opt/lagersystem
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm run start -- --hostname 0.0.0.0 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
